# supervisord.conf
[supervisord]
nodaemon=true

[program:dnsd]
command=/srv/bootstrapper-dnsd.sh
