# supervisord.conf
[supervisord]
nodaemon=true

[program:dnsd]
command=/srv/bootstrapper-dnsd.sh

[program:amqpd]
command=/srv/bootstrapper-amqpd.sh
